# deploy/orchard-deploy.service
[Unit]
Description=Pull latest orchard.rotko.net static site from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Pre-create the workdir so ProtectSystem/ReadWritePaths can bind-mount
# it on the first run; otherwise systemd fails at NAMESPACE setup.
ExecStartPre=/bin/mkdir -p /tmp/orchard-deploy
ExecStart=/usr/local/bin/orchard-deploy-install
SuccessExitStatus=0
User=root
ProtectSystem=full
ReadWritePaths=/var/www/orchard-rotko-net /tmp/orchard-deploy

[Install]
WantedBy=multi-user.target
