# deploy/orchard-deploy.service
[Unit]
Description=Pull latest orchard.rotko.net static site from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/orchard-deploy-install
SuccessExitStatus=0
User=root
# No sandbox: ProtectSystem + ReadWritePaths bind-mounts the docroot,
# which makes it "busy" so install.sh's atomic mv-swap of the dir
# can't succeed. Script is root-only on the edge boxes anyway.

[Install]
WantedBy=multi-user.target
